# data/etc/rc.d/init.d/shell
#!/bin/sh
### BEGIN INIT INFO
# Provides:            shell
# Required-Start:
# Should-Start:
# Required-Stop:
# Should-Stop:
# Default-Start:       3 4 5
# Default-Stop:
# Short-Description:   Runs emergency shell
# Description:         Runs emergency shell
### END INIT INFO

[ "$1" = start ] || exit 0

. /etc/init.d/template

start_daemon --background --lockfile "$LOCKFILE" --displayname "shell" -- setsid /bin/rdshell-forever
