# data/etc/rc.d/init.d/shell
#!/bin/sh
### BEGIN INIT INFO
# Provides:            shell
# Required-Start:      cmdline
# Should-Start:
# Required-Stop:
# Should-Stop:
# Default-Start:       3 4 5
# Default-Stop:
# Short-Description:   Prepare emergency shell.
# Description:         Prepare emergency shell.
### END INIT INFO

[ "$1" = start ] || exit 0

for f in /lib/shell/*; do
	[ -x "$f" ] || break
	"$f" &
done
