# features/nfsroot/data/etc/udev/rules.d/99-nfsroot.rules
SUBSYSTEM=="net", ACTION=="online", RUN+="/lib/udev-helpers/nfsroot"
