# testing/tests/ts0026-sort-services/data02
#!/bin/bash
### BEGIN INIT INFO
# Provides: service2
# Required-Start: service1
# Default-Start: 3
### END INIT INFO
