# testing/tests/ts0026-sort-services/data03
#!/bin/bash
### BEGIN INIT INFO
# Provides: service3
# Required-Start: service1
# Default-Start: 3
### END INIT INFO
