# features/btrfs/data/etc/udev/rules.d/99-initrd-btrfs.rules
SUBSYSTEM=="block", ACTION!="remove", ENV{ID_FS_TYPE}=="btrfs", RUN+="btrfs device scan $env{DEVNAME}"
