# testing/tests/ts0032-sort-services/rc.d/data10-sysctl
#!/bin/bash
### BEGIN INIT INFO
# Provides:            sysctl
# Required-Start:      mountvirtfs
# Default-Start:       3 4 5
### END INIT INFO
