# features/nvmf/data/etc/udev/rules.d/95-nvmf.rules
ACTION=="change", SUBSYSTEM=="fc", ENV{FC_EVENT}=="nvmediscovery", ENV{NVMEFC_HOST_TRADDR}=="*", ENV{NVMEFC_TRADDR}=="*", RUN+="/lib/uevent/filters/nvmf-fc"
