# testing/unit/sort-services/ts0014/data02
#!/bin/bash
### BEGIN INIT INFO
# Provides: service2
# Default-Start: 3
### END INIT INFO
