# features/raid/data/etc/udev/rules.d/99-initrd-raid.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/lib/udev-helpers/raid"
