# features/raid/data/etc/udev/rules.d/99-initrd-raid-md_run.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/lib/udev-helpers/md_run"
