# features/raid/data/etc/udev/rules.d/99-initrd-raid-md_run.rules
SUBSYSTEM=="block", ACTION!="remove", RUN+="/lib/uevent/filters/md_run"
