# features/mdadm/data/etc/udev/rules.d/64-md-raid-member.rules
SUBSYSTEM=="block", ACTION=="add", ENV{ID_FS_TYPE}=="linux_raid_member", RUN+="/lib/uevent/filters/md-raid-member"
