# testing/unit/sort-services/ts0009/data01
#!/bin/bash
### BEGIN INIT INFO
# Provides:            service1
# Required-Start:
# Should-Start:
# Required-Stop:
# Should-Stop:
# Default-Start:       3 4 5
# Default-Stop:
### END INIT INFO
