# testing/unit/sort-services/ts0011/data01
#!/bin/bash
### BEGIN INIT INFO
# Provides: service1
# Default-Start: 3
### END INIT INFO
