# testing/unit/sort-services/ts0016/rc.d/data06-plymouth
#!/bin/bash
### BEGIN INIT INFO
# Provides:            plymouth
# Required-Start:      udev
# X-Start-Before:      uevent
# Default-Start:       3 4 5
### END INIT INFO
