# testing/unit/sort-services/ts0016/rc.d/data08-runinit
#!/bin/bash
### BEGIN INIT INFO
# Provides:            runinit
# Should-Start:        killall
# Default-Start:       2
### END INIT INFO
