# testing/unit/sort-services/ts0016/rc.d/data12-uevent
#!/bin/bash
### BEGIN INIT INFO
# Provides:            uevent
# Required-Start:      cmdline fstab mountvirtfs
# Should-Start:        $syslog
# Default-Start:       3 4 5
### END INIT INFO
