# data/etc/udev/rules.d/00-debug.rules
SUBSYSTEM=="block", RUN+="/lib/uevent/filters/debug"
