# data/etc/udev/rules.d/80-network.rules
SUBSYSTEM=="net", ACTION=="add", RUN+="/lib/uevent/filters/network"
SUBSYSTEM=="net", ACTION=="online", RUN+="/lib/uevent/filters/mountdev"
