# data/etc/udev/rules.d/99-mountdev.rules
SUBSYSTEM=="block", ACTION!="remove", ENV{ID_FS_USAGE}=="filesystem", RUN+="/lib/uevent/filters/mountdev"
