# data/etc/udev/rules.d/99-netdev.rules
SUBSYSTEM=="net", ACTION=="add", RUN+="/lib/filters/netdev"
