# data/etc/udev/rules.d/99-netdev.rules
SUBSYSTEM=="net", ACTION=="online", RUN+="/lib/uevent/filters/mountdev"
SUBSYSTEM=="net", ACTION!="remove", RUN+="/lib/uevent/filters/netdev"
