# data/etc/udev/rules.d/99-netdev.rules
SUBSYSTEM=="net", ACTION!="remove", RUN+="/lib/filters/netdev"
