# data/etc/udev/rules.d/99-rootdev.rules
SUBSYSTEM=="block", ACTION=="add", ENV{ID_FS_USAGE}=="filesystem", RUN+="/lib/filters/rootdev"
