# data/etc/udev/rules.d/99-resume.rules
SUBSYSTEM=="block", ACTION!="remove", RUN+="/lib/initrd/filters/resume"
