# data/etc/udev/rules.d/99-resume.rules
SUBSYSTEM=="block", ACTION=="add", ENV{ID_FS_USAGE}=="filesystem", RUN+="/lib/udev-helpers/resume add"
