# testing/tests/ts0030-sort-services/data01
#!/bin/bash
### BEGIN INIT INFO
# Provides: service1 service2
# Required-Start: service2
# Default-Start: 3
### END INIT INFO
