# testing/units/sort-services/ts0015/data03
#!/bin/bash
### BEGIN INIT INFO
# Provides: service1 service2
# X-Start-Before: service2
# Default-Start: 3
### END INIT INFO
