# testing/units/sort-services/ts0015/data04
#!/bin/bash
### BEGIN INIT INFO
# Provides: service1 service2
# Required-Stop: service2
# Default-Stop: 3
### END INIT INFO
