# features/luks/data/etc/udev/rules.d/99-initrd-luks.rules
SUBSYSTEM=="block", ACTION!="remove", RUN+="/lib/filters/luks"
