# features/luks/data/etc/udev/rules.d/99-initrd-luks.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/lib/filters/luks"
