# features/lvm/data/etc/udev/rules.d/99-initrd-lvm.rules
SUBSYSTEM=="block", ACTION=="add", ENV{ID_FS_TYPE}=="LVM2_member", RUN+="/lib/filters/lvm"
