# features/lvm/data/etc/udev/rules.d/99-initrd-lvm.rules
SUBSYSTEM=="block", ACTION=="add", RUN+="/lib/udev-helpers/lvm add"
