# testing/units/sort-services/ts0016/rc.d/data01-fstab
#!/bin/bash
### BEGIN INIT INFO
# Provides:            fstab
# Required-Start:      cmdline
# Default-Start:       3 4 5
### END INIT INFO
