# testing/units/sort-services/ts0016/rc.d/data04-mountfs
#!/bin/bash
### BEGIN INIT INFO
# Provides:            mountfs mountvirtfs $local_fs
# Required-Start:      fstab
# Default-Start:       3 4 5
### END INIT INFO
