# testing/units/sort-services/ts0017/rc.d/data02-killall
#!/bin/bash
### BEGIN INIT INFO
# Provides:            killall
# Default-Start:       0 2 6
### END INIT INFO
