# testing/units/sort-services/ts0017/rc.d/data07-poll
#!/bin/bash
### BEGIN INIT INFO
# Provides:            poll
# Required-Start:      cmdline mountvirtfs
# Default-Start:       3 4 5
### END INIT INFO
