# testing/units/sort-services/ts0017/rc.d/data00-cmdline
#!/bin/bash
### BEGIN INIT INFO
# Provides:            cmdline
# Required-Start:      shell
# Default-Start:       3 4 5
### END INIT INFO
