# testing/units/sort-services/ts0017/rc.d/data05-mountsrc
#!/bin/bash
### BEGIN INIT INFO
# Provides:            mountsrc
# Required-Start:      udev uevent
# Default-Start:       3 4 5
### END INIT INFO
