# testing/units/sort-services/ts0017/rc.d/data09-shell
#!/bin/bash
### BEGIN INIT INFO
# Provides:            shell
# Default-Start:       3 4 5
### END INIT INFO
