# features/network/data/etc/udev/rules.d/81-network.rules
SUBSYSTEM=="net", RUN+="/lib/uevent/filters/network"
