# features/runtime/data/etc/udev/rules.d/80-network.rules
SUBSYSTEM=="net", ACTION=="online", RUN+="/lib/uevent/filters/mountdev"
