# testing/tests/ts0031-sort-services/rc.d/data03-modules
#!/bin/bash
### BEGIN INIT INFO
# Provides:            modules
# Required-Start:      sysctl cmdline
# Default-Start:       3 4 5
### END INIT INFO
