# features/kickstart/data/etc/rc.d/init.d/kickstart
#!/bin/bash
### BEGIN INIT INFO
# Provides:            kickstart
# Required-Start:      network-up udev uevent
# Should-Start:
# Required-Stop:
# Should-Stop:
# Default-Start:       3 4 5
# Default-Stop:
# Short-Description:   Runs kickstart script
# Description:         Runs kickstart script
# X-LFS-Provided-By:   LFS
### END INIT INFO

[ "$1" = start ] || exit 0

. /.initrd/initenv

[ -n "${KSFILE-}" ] || exit 0

[ -z "${KSFILE##/*}" ] ||
	KSFILE="/etc/ks.conf.d/$KSFILE"

if [ -s "$KSFILE" ]; then
	. uevent-sh-functions
	. initrd-sh-functions

	mkdir -p -- "$uevent_confdb/queue/pause/poll"

	kickstart -v "$KSFILE"

	rootdelay_reset_timer

	rmdir -- "$uevent_confdb/queue/pause/poll"
	rmdir -- "$uevent_confdb/queue/pause/udev"
	rmdir -- "$uevent_confdb/queue/pause/md-raid-member"
fi
